# src/root/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash
set -e

mkdir -p /run
chown abc:abc /run

exec php-fpm82 -F
